# supervisor.conf
[program:relays]
command=/usr/local/bin/python3.7 manage.py runserver 0:8000 
directory=/home/pi/django/mysite
redirect_stderr=true
stopsignal=INT
